# roles/tas_single_node/templates/systemd/backfill_redis.timer.j2
[Unit]
Description=Runs The backfill redis job

[Timer]
OnCalendar={{ tas_single_node_backfill_redis.schedule }}
Persistent=true
Unit=backfill_redis.service

[Install]
WantedBy=timers.target
