# roles/tas_single_node/templates/systemd/backup.service.j2
[Unit]
Description=Backup Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tas_backup.sh

[Install]
WantedBy]multi-user.target
