# roles/tas_single_node/templates/systemd/backfill_redis.timer.j2
[Unit]
Description=Runs The backfill redis job

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true
Unit=backfill_redis.service

[Install]
WantedBy=timers.target
