# roles/tas_single_node/templates/systemd/backup.timer.j2
[Unit]
Description= Run Backup Oneshot Job

[Timer]
OnCalendar={{ tas_single_node_backup_restore.backup.schedule }}
Persistent=true
Unit=tas_backup.service

[Install]
WantedBy=timers.target
