# roles/tas_single_node/templates/systemd/backup.service.j2
[Unit]
Description=Backup Service

[Service]
Type=oneshot
ExecStart=/bin/bash -c "/usr/local/bin/tas_backup.sh"
KillMode=process

[Install]
WantedBy=multi-user.target
